# debian/mdt-services.mdt-keymaster.service
[Unit]
Description=MDT public key management service
After=network.target

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/bin/mdt-keymaster
Restart=on-failure
User=mendel
Group=mendel
WorkingDirectory=~

[Install]
WantedBy=multi-user.target
